# etc/update-motd.d/00-moodeos-header
#!/bin/bash

REL=$(awk '/Release: /{print $2 " " $3;}' /var/www/footer.php | tr -d ',')
echo Kernel $(uname -r) $(uname -v | awk '{print $1}'), Raspbian $(cat /etc/debian_version), moOde $REL
