# etc/update-motd.d/00-moodeos-header
#!/bin/bash

MOODE_RELEASE=$(moodeutl --mooderel | tr -d '\n')
HDWR_REV=$(sqlite3 /var/local/www/db/moode-sqlite3.db "SELECT value FROM cfg_system WHERE param='hdwrrev'")
[ $(uname -m) = "aarch64" ] && KERNEL_ARCH="64-bit" || KERNEL_ARCH="32-bit"
AUDIO_DEV=$(sqlite3 /var/local/www/db/moode-sqlite3.db "SELECT value FROM cfg_system WHERE param='adevname'")

echo Linux: $(uname -vr | awk '{print $1" "$2}') \| $(uname -m) \($KERNEL_ARCH\) \| "RaspiOS: "$(cat /etc/debian_version)
echo Model: $HDWR_REV
echo Audio: $AUDIO_DEV

echo "

                      ____     __      ___
          __ _  ___  / __ \___/ /__   ( _ )
         /  ' \/ _ \/ /_/ / _  / -_) / _  |
        /_/_/_/\___/\____/\_,_/\__/  \___/

                 moOde audio player
              Release $MOODE_RELEASE
                (C) 2014 Tim Curtis


The programs included with moOde are free software; the
exact distribution terms for each program are described
in the individual files in /usr/share/doc/*/copyright or
in the moOde source files. moOde comes with absolutely
no warranties either expressed or implied, or any other
such guarantees.

To set a new password for your userid type 'passwd'.
"
